# root/etc/services.d/pritunl/run
#!/usr/bin/with-contenv bash

echo "Starting pritunl"
exec /usr/bin/pritunl start
